# service/com.dodi.rotate-logs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dodi.rotate-logs</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/mokie/services/scripts/rotate-logs.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/mokie</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/mokie/services/scripts/logs/rotate-logs.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mokie/services/scripts/logs/rotate-logs.log</string>
</dict>
</plist>
